# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/local/bin/ovpn_run

[program:socat]
startsecs=0
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=socat %(ENV_SOCAT_SOURCE)s %(ENV_SOCAT_DEST)s